# cross/lxc/PLIST
bin:bin/lxc-attach
bin:bin/lxc-autostart
bin:bin/lxc-cgroup
bin:bin/lxc-checkconfig
bin:bin/lxc-checkpoint
bin:bin/lxc-config
bin:bin/lxc-console
bin:bin/lxc-copy
bin:bin/lxc-create
bin:bin/lxc-destroy
bin:bin/lxc-device
bin:bin/lxc-execute
bin:bin/lxc-freeze
bin:bin/lxc-info
bin:bin/lxc-ls
bin:bin/lxc-monitor
bin:bin/lxc-snapshot
bin:bin/lxc-start
bin:bin/lxc-stop
bin:bin/lxc-top
bin:bin/lxc-unfreeze
bin:bin/lxc-unshare
bin:bin/lxc-update-config
bin:bin/lxc-usernsexec
bin:bin/lxc-wait
rsc:etc/apparmor.d/abstractions/lxc/container-base
rsc:etc/apparmor.d/abstractions/lxc/start-container
rsc:etc/apparmor.d/lxc-containers
rsc:etc/apparmor.d/lxc/lxc-default
rsc:etc/apparmor.d/lxc/lxc-default-cgns
rsc:etc/apparmor.d/lxc/lxc-default-with-mounting
rsc:etc/apparmor.d/lxc/lxc-default-with-nesting
rsc:etc/apparmor.d/usr.bin.lxc-start
rsc:etc/bash_completion.d/lxc
rsc:etc/default/lxc
rsc:etc/default/lxc-net
rsc:etc/init.d/lxc
rsc:etc/init.d/lxc-net
rsc:etc/lxc/default.conf
lnk:lib/liblxc.so
lnk:lib/liblxc.so.1
lib:lib/liblxc.so.1.7.0
rsc:lib/lxc/rootfs/README
rsc:libexec/lxc/hooks/unmount-namespace
rsc:libexec/lxc/lxc-apparmor-load
rsc:libexec/lxc/lxc-containers
bin:libexec/lxc/lxc-monitord
rsc:libexec/lxc/lxc-net
bin:libexec/lxc/lxc-user-nic
bin:sbin/init.lxc
rsc:share/doc/lxc/examples/lxc-complex.conf
rsc:share/doc/lxc/examples/lxc-empty-netns.conf
rsc:share/doc/lxc/examples/lxc-macvlan.conf
rsc:share/doc/lxc/examples/lxc-no-netns.conf
rsc:share/doc/lxc/examples/lxc-phys.conf
rsc:share/doc/lxc/examples/lxc-veth.conf
rsc:share/doc/lxc/examples/lxc-vlan.conf
rsc:share/doc/lxc/examples/seccomp-v1.conf
rsc:share/doc/lxc/examples/seccomp-v2-denylist.conf
rsc:share/doc/lxc/examples/seccomp-v2.conf
rsc:share/lxc/config/common.conf
rsc:share/lxc/config/common.conf.d/README
rsc:share/lxc/config/common.seccomp
rsc:share/lxc/config/nesting.conf
rsc:share/lxc/config/oci.common.conf
rsc:share/lxc/config/userns.conf
rsc:share/lxc/hooks/clonehostname
rsc:share/lxc/hooks/dhclient
rsc:share/lxc/hooks/dhclient-script
rsc:share/lxc/hooks/mountecryptfsroot
rsc:share/lxc/hooks/nvidia
rsc:share/lxc/hooks/squid-deb-proxy-client
rsc:share/lxc/hooks/ubuntu-cloud-prep
rsc:share/lxc/lxc-patch.py
rsc:share/lxc/lxc.functions
rsc:share/lxc/selinux/lxc.if
rsc:share/lxc/selinux/lxc.te
rsc:share/lxc/templates/lxc-busybox
rsc:share/lxc/templates/lxc-download
rsc:share/lxc/templates/lxc-local
rsc:share/lxc/templates/lxc-oci
