# spk/safekeep/src/safekeep.cron
MAILTO=""
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/syno/sbin:/usr/syno/bin:/usr/local/sbin:/usr/local/bin
#minute hour mday month wday who command
00 01 * * * root /usr/local/safekeep/bin/safekeep-backups.sh
