# cross/libpam/PLIST
lib/systemd/system/pam_namespace.service
share/doc/Linux-PAM/draft-morgan-pam-current.txt
share/doc/Linux-PAM/index.html
share/doc/Linux-PAM/rfc86.0.txt
share/locale/af/LC_MESSAGES/Linux-PAM.mo
share/locale/am/LC_MESSAGES/Linux-PAM.mo
share/locale/ar/LC_MESSAGES/Linux-PAM.mo
share/locale/as/LC_MESSAGES/Linux-PAM.mo
share/locale/az/LC_MESSAGES/Linux-PAM.mo
share/locale/be/LC_MESSAGES/Linux-PAM.mo
share/locale/bg/LC_MESSAGES/Linux-PAM.mo
share/locale/bn/LC_MESSAGES/Linux-PAM.mo
share/locale/bn_IN/LC_MESSAGES/Linux-PAM.mo
share/locale/bs/LC_MESSAGES/Linux-PAM.mo
share/locale/ca/LC_MESSAGES/Linux-PAM.mo
share/locale/cs/LC_MESSAGES/Linux-PAM.mo
share/locale/cy/LC_MESSAGES/Linux-PAM.mo
share/locale/da/LC_MESSAGES/Linux-PAM.mo
share/locale/de/LC_MESSAGES/Linux-PAM.mo
share/locale/de_CH/LC_MESSAGES/Linux-PAM.mo
share/locale/el/LC_MESSAGES/Linux-PAM.mo
share/locale/eo/LC_MESSAGES/Linux-PAM.mo
share/locale/es/LC_MESSAGES/Linux-PAM.mo
share/locale/et/LC_MESSAGES/Linux-PAM.mo
share/locale/eu/LC_MESSAGES/Linux-PAM.mo
share/locale/fa/LC_MESSAGES/Linux-PAM.mo
share/locale/fi/LC_MESSAGES/Linux-PAM.mo
share/locale/fr/LC_MESSAGES/Linux-PAM.mo
share/locale/ga/LC_MESSAGES/Linux-PAM.mo
share/locale/gl/LC_MESSAGES/Linux-PAM.mo
share/locale/gu/LC_MESSAGES/Linux-PAM.mo
share/locale/he/LC_MESSAGES/Linux-PAM.mo
share/locale/hi/LC_MESSAGES/Linux-PAM.mo
share/locale/hr/LC_MESSAGES/Linux-PAM.mo
share/locale/hu/LC_MESSAGES/Linux-PAM.mo
share/locale/ia/LC_MESSAGES/Linux-PAM.mo
share/locale/id/LC_MESSAGES/Linux-PAM.mo
share/locale/is/LC_MESSAGES/Linux-PAM.mo
share/locale/it/LC_MESSAGES/Linux-PAM.mo
share/locale/ja/LC_MESSAGES/Linux-PAM.mo
share/locale/ka/LC_MESSAGES/Linux-PAM.mo
share/locale/kk/LC_MESSAGES/Linux-PAM.mo
share/locale/km/LC_MESSAGES/Linux-PAM.mo
share/locale/kn/LC_MESSAGES/Linux-PAM.mo
share/locale/ko/LC_MESSAGES/Linux-PAM.mo
share/locale/kw_GB/LC_MESSAGES/Linux-PAM.mo
share/locale/ky/LC_MESSAGES/Linux-PAM.mo
share/locale/lt/LC_MESSAGES/Linux-PAM.mo
share/locale/lv/LC_MESSAGES/Linux-PAM.mo
share/locale/mk/LC_MESSAGES/Linux-PAM.mo
share/locale/ml/LC_MESSAGES/Linux-PAM.mo
share/locale/mn/LC_MESSAGES/Linux-PAM.mo
share/locale/mr/LC_MESSAGES/Linux-PAM.mo
share/locale/ms/LC_MESSAGES/Linux-PAM.mo
share/locale/my/LC_MESSAGES/Linux-PAM.mo
share/locale/nb/LC_MESSAGES/Linux-PAM.mo
share/locale/ne/LC_MESSAGES/Linux-PAM.mo
share/locale/nl/LC_MESSAGES/Linux-PAM.mo
share/locale/nn/LC_MESSAGES/Linux-PAM.mo
share/locale/or/LC_MESSAGES/Linux-PAM.mo
share/locale/pa/LC_MESSAGES/Linux-PAM.mo
share/locale/pl/LC_MESSAGES/Linux-PAM.mo
share/locale/pt/LC_MESSAGES/Linux-PAM.mo
share/locale/pt_BR/LC_MESSAGES/Linux-PAM.mo
share/locale/ro/LC_MESSAGES/Linux-PAM.mo
share/locale/ru/LC_MESSAGES/Linux-PAM.mo
share/locale/si/LC_MESSAGES/Linux-PAM.mo
share/locale/sk/LC_MESSAGES/Linux-PAM.mo
share/locale/sl/LC_MESSAGES/Linux-PAM.mo
share/locale/sq/LC_MESSAGES/Linux-PAM.mo
share/locale/sr/LC_MESSAGES/Linux-PAM.mo
share/locale/sr@latin/LC_MESSAGES/Linux-PAM.mo
share/locale/sv/LC_MESSAGES/Linux-PAM.mo
share/locale/ta/LC_MESSAGES/Linux-PAM.mo
share/locale/te/LC_MESSAGES/Linux-PAM.mo
share/locale/tg/LC_MESSAGES/Linux-PAM.mo
share/locale/th/LC_MESSAGES/Linux-PAM.mo
share/locale/tr/LC_MESSAGES/Linux-PAM.mo
share/locale/uk/LC_MESSAGES/Linux-PAM.mo
share/locale/ur/LC_MESSAGES/Linux-PAM.mo
share/locale/vi/LC_MESSAGES/Linux-PAM.mo
share/locale/yo/LC_MESSAGES/Linux-PAM.mo
share/locale/zh_CN/LC_MESSAGES/Linux-PAM.mo
share/locale/zh_HK/LC_MESSAGES/Linux-PAM.mo
share/locale/zh_TW/LC_MESSAGES/Linux-PAM.mo
share/locale/zu/LC_MESSAGES/Linux-PAM.mo
share/man/man3/misc_conv.3
share/man/man3/pam.3
share/man/man3/pam_acct_mgmt.3
share/man/man3/pam_authenticate.3
share/man/man3/pam_chauthtok.3
share/man/man3/pam_close_session.3
share/man/man3/pam_conv.3
share/man/man3/pam_end.3
share/man/man3/pam_error.3
share/man/man3/pam_fail_delay.3
share/man/man3/pam_get_authtok.3
share/man/man3/pam_get_authtok_noverify.3
share/man/man3/pam_get_authtok_verify.3
share/man/man3/pam_get_data.3
share/man/man3/pam_get_item.3
share/man/man3/pam_get_user.3
share/man/man3/pam_getenv.3
share/man/man3/pam_getenvlist.3
share/man/man3/pam_info.3
share/man/man3/pam_misc_drop_env.3
share/man/man3/pam_misc_paste_env.3
share/man/man3/pam_misc_setenv.3
share/man/man3/pam_open_session.3
share/man/man3/pam_prompt.3
share/man/man3/pam_putenv.3
share/man/man3/pam_set_data.3
share/man/man3/pam_set_item.3
share/man/man3/pam_setcred.3
share/man/man3/pam_sm_acct_mgmt.3
share/man/man3/pam_sm_authenticate.3
share/man/man3/pam_sm_chauthtok.3
share/man/man3/pam_sm_close_session.3
share/man/man3/pam_sm_open_session.3
share/man/man3/pam_sm_setcred.3
share/man/man3/pam_start.3
share/man/man3/pam_strerror.3
share/man/man3/pam_syslog.3
share/man/man3/pam_verror.3
share/man/man3/pam_vinfo.3
share/man/man3/pam_vprompt.3
share/man/man3/pam_vsyslog.3
share/man/man3/pam_xauth_data.3
share/man/man5/access.conf.5
share/man/man5/environment.5
share/man/man5/faillock.conf.5
share/man/man5/group.conf.5
share/man/man5/limits.conf.5
share/man/man5/namespace.conf.5
share/man/man5/pam.conf.5
share/man/man5/pam.d.5
share/man/man5/pam_env.conf.5
share/man/man5/time.conf.5
share/man/man8/PAM.8
share/man/man8/faillock.8
share/man/man8/mkhomedir_helper.8
share/man/man8/pam.8
share/man/man8/pam_access.8
share/man/man8/pam_debug.8
share/man/man8/pam_deny.8
share/man/man8/pam_echo.8
share/man/man8/pam_env.8
share/man/man8/pam_exec.8
share/man/man8/pam_faildelay.8
share/man/man8/pam_faillock.8
share/man/man8/pam_filter.8
share/man/man8/pam_ftp.8
share/man/man8/pam_group.8
share/man/man8/pam_issue.8
share/man/man8/pam_keyinit.8
share/man/man8/pam_lastlog.8
share/man/man8/pam_limits.8
share/man/man8/pam_listfile.8
share/man/man8/pam_localuser.8
share/man/man8/pam_loginuid.8
share/man/man8/pam_mail.8
share/man/man8/pam_mkhomedir.8
share/man/man8/pam_motd.8
share/man/man8/pam_namespace.8
share/man/man8/pam_namespace_helper.8
share/man/man8/pam_nologin.8
share/man/man8/pam_permit.8
share/man/man8/pam_pwhistory.8
share/man/man8/pam_rhosts.8
share/man/man8/pam_rootok.8
share/man/man8/pam_securetty.8
share/man/man8/pam_setquota.8
share/man/man8/pam_shells.8
share/man/man8/pam_stress.8
share/man/man8/pam_succeed_if.8
share/man/man8/pam_time.8
share/man/man8/pam_timestamp.8
share/man/man8/pam_timestamp_check.8
share/man/man8/pam_umask.8
share/man/man8/pam_unix.8
share/man/man8/pam_usertype.8
share/man/man8/pam_warn.8
share/man/man8/pam_wheel.8
share/man/man8/pam_xauth.8
share/man/man8/pwhistory_helper.8
share/man/man8/unix_chkpwd.8
share/man/man8/unix_update.8
