# cross/mesa/PLIST
lib:lib/libvulkan_intel.so
rsc:etc/vulkan/icd.d
rsc:share/drirc.d/00-mesa-defaults.conf
